# dist/c-sync.service
[Unit]
Description=C-Sync Atualizador Automático Digitro
After=syslog.target network.target

[Service]
Type=forking
User=root
Group=root

ExecStart=/home2/c-sync-server/init.d/c-sync-startup.sh start
ExecStop=/home2/c-sync-server/init.d/c-sync-startup.sh stop
ExecRestart=/home2/c-sync-server/init.d/c-sync-startup.sh restart
PIDFile=/home2/c-sync-server/pid/c-sync.pid

[Install]
WantedBy=multi-user.target